# ft_services/srcs/influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:influx]
command = sh -c "influxd run -config /etc/influxdb.conf; kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "cd /telegraf-1.15.2/usr/bin; ./telegraf --config /telegraf-1.15.2/etc/telegraf/telegraf.conf; kill -s SIGTERM $(cat supervisord.pid)"
