# ft_services/srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command = sh -c "/usr/bin/mysqld_safe; kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "cd /telegraf-1.15.2/usr/bin; ./telegraf --config /telegraf-1.15.2/etc/telegraf/telegraf.conf; kill -s SIGTERM $(cat supervisord.pid)"