# ft_services/srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command = sh -c "vsftpd /etc/vsftpd/vsftpd.conf; kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "cd /telegraf-1.15.2/usr/bin; ./telegraf --config /telegraf-1.15.2/etc/telegraf/telegraf.conf; kill -s SIGTERM $(cat supervisord.pid)"